# supervisord.conf
[supervisord]
nodaemon = true
pidfile = /tmp/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[include]
files = /etc/supervisord/conf.d/*.conf

[program:web]
command=uwsgi --show-config
stdout_events_enabled=true
stderr_events_enabled=true

[program:bot]
command=python manage.py run_bot
stdout_events_enabled=true
stderr_events_enabled=true

[eventlistenter:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
